# ansible/roles/openclaw/templates/openclaw.service.j2
[Unit]
Description=OpenClaw AI Agent Gateway
Documentation=https://github.com/openclaw/openclaw
After=network-online.target docker.service
{% if tailscale_enabled | default(false) %}
After=tailscaled.service
{% endif %}
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User={{ openclaw_user }}
Group={{ openclaw_user }}
WorkingDirectory={{ openclaw_home }}

# Environment — aligned with upstream pnpm paths
Environment="NODE_ENV=production"
Environment="PNPM_HOME={{ openclaw_home }}/.local/share/pnpm"
Environment="PATH={{ openclaw_home }}/.local/bin:{{ openclaw_home }}/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin"
Environment="HOME={{ openclaw_home }}"
{% if openclaw_uid_value is defined %}
Environment="XDG_RUNTIME_DIR=/run/user/{{ openclaw_uid_value }}"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{ openclaw_uid_value }}/bus"
{% endif %}
EnvironmentFile=-{{ openclaw_config_dir }}/.env

# Start OpenClaw (upstream uses `openclaw gateway` as the entry point)
ExecStart={{ openclaw_home }}/.local/bin/openclaw gateway

# Restart policy
Restart=always
RestartSec=10

# Security hardening (aligned with upstream openclaw-host.service.j2)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ openclaw_config_dir }}
ReadWritePaths={{ openclaw_home }}/.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
